# judgehost/conf/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/var/run/supervisord.pid

[program:judgedaemon]
command=judgedaemon
